# Buildroot/external/external_config_systemd/package/hello-systemd/hello-systemd.service
[Unit]
Description=GNU Hello World startup script
After=multi-user.target

# Define the command to start here
[Service]
ExecStart=/usr/bin/hello-systemd.sh
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target

